# scripts/supervisor/onwynd.conf
; =============================================================================
; Onwynd — Supervisor configuration
; Manages queue workers, scheduler, and Reverb (WebSocket) for production.
;
; Install:
;   sudo cp onwynd.conf /etc/supervisor/conf.d/onwynd.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start onwynd-worker:*
;
; Requires: supervisord running on the host
;           PHP path set correctly (update /usr/bin/php if needed)
;           BASE_DIR set to the actual Laravel root
; =============================================================================

[group:onwynd-worker]
programs=onwynd-queue-default,onwynd-queue-notifications,onwynd-queue-ai,onwynd-queue-payments,onwynd-scheduler

; ---------------------------------------------------------------------------
; Default queue — general jobs (emails, webhooks, misc)
; ---------------------------------------------------------------------------
[program:onwynd-queue-default]
command=/usr/bin/php /var/www/onwynd/api/artisan queue:work redis --queue=default,mail,webhooks --sleep=3 --tries=3 --max-time=3600 --memory=128
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=600
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/onwynd-queue-default.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=www-data
environment=HOME="/var/www",USER="www-data"

; ---------------------------------------------------------------------------
; Notifications queue — session reminders, push, SMS
; ---------------------------------------------------------------------------
[program:onwynd-queue-notifications]
command=/usr/bin/php /var/www/onwynd/api/artisan queue:work redis --queue=notifications --sleep=3 --tries=3 --max-time=3600 --memory=128
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=600
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/onwynd-queue-notifications.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=www-data
environment=HOME="/var/www",USER="www-data"

; ---------------------------------------------------------------------------
; AI queue — OpenAI/Groq calls, document analysis
; ---------------------------------------------------------------------------
[program:onwynd-queue-ai]
command=/usr/bin/php /var/www/onwynd/api/artisan queue:work redis --queue=ai --sleep=3 --tries=2 --max-time=3600 --memory=256
process_name=%(program_name)s_%(process_num)02d
numprocs=3
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=600
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/onwynd-queue-ai.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=www-data
environment=HOME="/var/www",USER="www-data"

; ---------------------------------------------------------------------------
; Payments queue — Stripe/Paystack webhook processing
; ---------------------------------------------------------------------------
[program:onwynd-queue-payments]
command=/usr/bin/php /var/www/onwynd/api/artisan queue:work redis --queue=payments --sleep=3 --tries=3 --max-time=3600 --memory=128
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=600
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/onwynd-queue-payments.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=www-data
environment=HOME="/var/www",USER="www-data"

; ---------------------------------------------------------------------------
; Laravel Scheduler — runs `artisan schedule:run` every minute
; ---------------------------------------------------------------------------
[program:onwynd-scheduler]
command=/usr/bin/php /var/www/onwynd/api/artisan schedule:work
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/onwynd-scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
user=www-data
environment=HOME="/var/www",USER="www-data"
